# deploy/systemd/corelasi-backup.timer
[Unit]
Description=Run the CORELASI backup every day

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
